# src/openntpd.service.in
[Unit]
Description=OpenNTPD Network Time Protocol
Documentation=man:openntpd(8)
After=network.target

[Service]
Type=forking
ExecStartPre=@sbindir@/ntpd -n $DAEMON_OPTS
ExecStart=@sbindir@/ntpd $DAEMON_OPTS
ExecReload=@sbindir@/ntpd -n $DAEMON_OPTS
KillMode=process
Restart=on-failure
EnvironmentFile=-@sysconfdir@/default/openntpd
RuntimeDirectory=openntpd
RuntimeDirectoryMode=0755
StateDirectory=openntpd
StateDirectoryMode=0755

[Install]
WantedBy=multi-user.target
